# task1_3/task1_jboss/templates/wildfly.service.erb
[Unit]
Description=WildFly 10.1.0.Final
After=network.target

[Service]
Type=simple
ExecStart=/opt/wildfly/bin/standalone.sh

[Install]
WantedBy=multi-user.target
